# rc.local
python3 /root/driveMonitor.py &

exit 0
